# archive/python/deploy/wg-manager-worker.service
# WireGuard Web Manager - 特権 Worker（root 実行）
# Web アプリは一般ユーザーで動作し、wg 操作のみこの Worker にソケット経由で依頼する。
#
# 配置手順（手元のターミナルで実行）:
#   sudo cp /home/kanri/wireguard-web-manager/deploy/wg-manager-worker.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now wg-manager-worker
#   sudo systemctl status wg-manager-worker
#
# プロジェクトを別ディレクトリに置いた場合は、WorkingDirectory / Environment / ExecStart のパスを編集すること。
[Unit]
Description=WireGuard Web Manager Worker (root)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/kanri/wireguard-web-manager
Environment=CONFIG_PATH=/home/kanri/wireguard-web-manager/config.yaml
ExecStart=/home/kanri/wireguard-web-manager/.venv/bin/python3 /home/kanri/wireguard-web-manager/wg_worker.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
