# deploy/wg-manager-worker.service
# WireGuard Web Manager - 特権 Worker（Rust 版・root 実行）
# 使い方: パスを環境に合わせて編集し、rust-wg-manager-worker.service として
#         /etc/systemd/system/ にコピーしてから daemon-reload / enable --now する。
[Unit]
Description=WireGuard Web Manager Worker (Rust - root)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/kanri/wireguard-web-manager
Environment=CONFIG_PATH=/home/kanri/wireguard-web-manager/config.yaml
ExecStart=/home/kanri/wireguard-web-manager/rust/target/release/wg-worker
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
