# deploy/wireguard-webmanager-worker.service
[Unit]
Description=WireGuard Web Manager - Privileged Worker
Documentation=https://github.com/tm-io/WireGuard_WebManager
After=network.target
# wg0 インターフェースが存在する環境では起動順序を保証できる
# After=sys-subsystem-net-devices-wg0.device

[Service]
Type=simple
User=root
Group=root
Environment=CONFIG_PATH=/etc/wireguard-webmanager/config.yaml
ExecStart=/usr/bin/wg-worker
Restart=on-failure
RestartSec=5
# ソケットファイルの親ディレクトリを事前に確保
RuntimeDirectory=wg-manager
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
